# files/lxc-backup@.service
[Unit]
Description=LXC container backup: %i

[Service]
Type=simple
ExecStart=/bin/bash /usr/local/bin/lxc-backup.sh %i
